# moneytolia-admin.service
[Unit]
Description=Paypango admin service.

[Service]
Type=simple
User=paypango-admin
ExecStart=/bin/bash /opt/monpaypangoeytolia/admin/run.sh

[Install]
WantedBy=multi-user.target

